# psiphon
 #! /bin/sh
### BEGIN INIT INFO
# Provides:          psiphon
# Required-Start:    $start
# Required-Stop:     $shutdown
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: psiphine init scipt
# Description:       init script for the psiphon proxy client
### END INIT INFO
# Author: Alexander Wiechert

# Aktionen
PSIPHON="psiphon-tunnel-core-x86_64 -config client.free.conf"
PROCESS=$(ps -ef | grep  "psiphon-tunnel-core-x86_64 -config client.free.conf" |grep -v grep | awk {'print $2'})
DATE=$(date +'%Y-%m-%d')

case $1 in

  start)
        $PSIPHON 2>&1 >> /var/log/psiphon/$DATE-psiphon.log &
    ;;

  stop) 
        
        kill -5 $PROCESS
    ;;

  status)
	if [ $PROCESS != 0 ]; then
	  echo "psiphon client is running"
	  tail -f /var/log/psiphon/$DATE-psiphon.log
	else 
	  echo "psiphon client is NOT running"

	fi	
    ;;

  *)
    echo "$0 start | stop | status"
    ;;
esac
